# etc/systemd/system/recipes.service
[Unit]
Description=uWSGI instance to serve Hiren-Recipes

[Service]
ExecStart=/bin/bash -c 'su prism; cd /home/prism/prod/Hiren-Recipes; source .env/bin/activate; uwsgi --enable-threads --ini hiren.ini'

[Install]
WantedBy=multi-user.target
